# src/etc/supervisord.conf
[supervisord]
nodaemon=true
user=root
pidfile=/run/supervisord.pid
logfile=/var/log/supervisord/supervisord.log

[program:nginx]
command=nginx -g "daemon off;"

[program:cron]
command=crond -f
